# hack/cloudflared-k8s-proxy.service
[Unit]
Description=Cloudflare Access proxy for k8s.thecluster.io
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=cloudflared access tcp --hostname k8s.thecluster.io --url localhost:16443
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
